# docker/app/crontab
* * * * * /usr/bin/php7.3 /var/www/app/artisan schedule:run >> /var/www/app/storage/logs/cron.log 2>&1
